# modules/homeassistant/templates/socat.service.erb
[Unit]
Description=socat <%= @name %>
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat pty,link=<%= @local_socket %>,raw,mode=666 tcp:<%= @remote_addr %>
