# modules/nginx_proxy/templates/sso-proxy.service.erb
[Unit]
Description=sso-proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sso-proxy
Environment=EMAIL_ADDRESSES=<%= @allowed_emails %>
Environment=UPSTREAM_CONFIGS=/etc/sso/upstream_configs.yml
Environment=PROVIDER_URL=https://auth.sso.boyfriend.network
Environment=CLIENT_ID=<%= @proxy_client_id %>
Environment=CLIENT_SECRET=<%= @proxy_client_secret %>
Environment=AUTH_CODE_SECRET=<%= @proxy_auth_code_secret %>
Environment=COOKIE_SECRET=<%= @proxy_cookie_secret %>
Environment=STATSD_HOST=127.0.0.1
Environment=STATSD_PORT=8125
Environment=CLUSTER=prod
Environment=PORT=<%= @proxy_port %>
