# modules/nginx_proxy/templates/sso-auth.service.erb
[Unit]
Description=sso-auth
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sso-auth
Environment=HOST=auth.sso.boyfriend.network
Environment=SSO_EMAIL_ADDRESSES=<%= @allowed_emails %>
Environment=REDIRECT_URL=https://auth.sso.boyfriend.network/oauth/callback
Environment=PROXY_ROOT_DOMAIN=boyfriend.network
Environment=CLIENT_ID=<%= @client_id %>
Environment=CLIENT_SECRET=<%= @client_secret %>
Environment=PROXY_CLIENT_ID=<%= @proxy_client_id %>
Environment=PROXY_CLIENT_SECRET=<%= @proxy_client_secret %>
Environment=AUTH_CODE_SECRET=<%= @auth_auth_code_secret %>
Environment=COOKIE_SECRET=<%= @auth_cookie_secret %>
Environment=STATSD_HOST=127.0.0.1
Environment=STATSD_PORT=8125
Environment=CLUSTER=prod
Environment=PORT=<%= @auth_port %>
